# .config/systemd/user/kctx-health.service
[Unit]
Description=Refresh kctx connection health
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
Environment=PATH=%h/dev/google-cloud-sdk/bin:%h/.local/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=%h/.local/bin/kctx health --all --refresh --jobs 4
Nice=10
IOSchedulingClass=idle
TimeoutStartSec=10min
StandardOutput=null
StandardError=journal
